# _system/com.bitcoinbeacon.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bitcoinbeacon.autopush</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mzahn/bitcoin beacon/_system/auto-push.sh</string>
  </array>
  <!-- Every 30 minutes. Runs once at load, then every 1800s. If the Mac was
       asleep at a scheduled time, launchd fires once on wake — so a missed
       morning push goes out shortly after the machine comes back. -->
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/bitcoinbeacon-push.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bitcoinbeacon-push.log</string>
</dict>
</plist>
